# systemd/bb-nrepl.socket
[Unit]
Description=Babashka nREPL Socket

[Socket]
ListenStream=7888
Accept=false
